# install/system/etc/profile.d/90-archmeros-build.sh
export MAKEFLAGS='-j2'
export CMAKE_BUILD_PARALLEL_LEVEL='2'
