# config/systemd/user/content-factory-forum-link-health.service
[Unit]
Description=Content Factory Forum Link Health Service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/zacmero/projects/content-factory-stack
Environment=FORUM_LINK_HEALTH_HOST=0.0.0.0
Environment=FORUM_LINK_HEALTH_PORT=8791
Environment=FORUM_LINK_HEALTH_BROWSER=chromium
ExecStart=/usr/bin/node /home/zacmero/projects/content-factory-stack/scripts/forum_link_health_service.mjs
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
